# pi/lyra-webapp.service
[Unit]
Description=Lyra Coffee Machine Next.js Webapp
After=network.target

[Service]
Type=simple
User=pi
WorkingDirectory=/opt/lyra/webapp-next
ExecStart=/usr/bin/node_modules/.bin/next start -p 3000
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
Environment=NODE_ENV=production
EnvironmentFile=/opt/lyra/webapp-next/.env.local

[Install]
WantedBy=multi-user.target
